# systemd/assist-volume-restore.service
[Unit]
Description=Restore Assist volume
After=sound.target
Wants=sound.target

[Service]
Type=oneshot
# Wait for ALSA to settle, then retry a few times
ExecStart=/bin/bash -lc 'for i in {1..10}; do /usr/bin/amixer -c 1 cset name="Speaker Playback Volume" 110 && exit 0; sleep 1; done; exit 1'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
