# systemd/ha-satellite-mq-agent.service
[Unit]
Description=HA Satellite Agent (MQTT)
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=root
Restart=always
RestartSec=2

# Single source of truth
EnvironmentFile=-/etc/ha-satellite/satellite.env

# Defaults (can be overridden in satellite.env)
Environment=MQTT_BASE_TOPIC=ha/satellite
Environment=MQTT_PORT=1883
Environment=MQTT_QOS=0
Environment=PUBLISH_STATE=1
Environment=RETAIN_STATE=0

# Volume / ALSA defaults (keep here or move into env if you want)
Environment=ALSA_CARD=1
Environment=ALSA_NUMID=13
Environment=VOL_MIN=0
Environment=VOL_MAX=255

# Fail fast if identity is missing (prevents anonymous topics / weird behavior)
ExecStartPre=/bin/sh -c 'test -n "$SAT_HOSTNAME"'

ExecStart=/usr/local/bin/ha-satellite-mq-agent.py

[Install]
WantedBy=multi-user.target
